# files/puppyos/etc/runit/sv/dhcpcd/run
#!/bin/sh
mkdir -p /var/db/dhcpcd
exec /usr/bin/dhcpcd -B
