# files/puppyos/etc/runit/sv/dash-ttyS0/run
#!/bin/sh
exec setsid -c dash -limE </dev/ttyS0 >/dev/ttyS0 2>&1
